# etc/supervisor/conf.d/pi_weather.conf
[program:pi_weather]
command=python /home/pi/git/PiTFTWeather/PiTFTWeather.py
user=root
directory=/home/pi/git/PiTFTWeather
autostart=false
autorestart=true
